# .config/systemd/user/addressbook.timer
[Unit]
Description=maildir-rank-addr update timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h

[Install]
WantedBy=default.target
