# .config/systemd/user/fastsync.timer
[Unit]
Description=fast mail sync timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1m

[Install]
WantedBy=default.target
