# .config/systemd/user/mailsync-low.timer
[Unit]
Description=Mail Sync low frequency timer

[Timer]
OnBootSec=1h
OnUnitActiveSec=6h

[Install]
WantedBy=default.target
