# .config/systemd/user/mailsync-medium.timer
[Unit]
Description=Mail Sync medium frequency timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15m

[Install]
WantedBy=default.target
