# .config/systemd/user/mbsync-run.timer
[Unit]
Description=Mbsync timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=2m

[Install]
WantedBy=default.target
