# .config/systemd/user/mbsync-run.timer
[Unit]
Description=Mbsync timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1m
Persistent=true

[Install]
WantedBy=default.target
