# .config/systemd/user/restic.timer
[Unit]
Description=Restic backup timer

[Timer]
OnCalendar=daily
Persistent=true


[Install]
WantedBy=default.target
