# .config/systemd/user/sync-org-git.timer
[Unit]
Description=Sync org git timer

[Timer]
OnBootSec=0min
OnCalendar=*:*:0/5

[Install]
WantedBy=default.target
