# .config/systemd/user/tasksync.timer
[Unit]
Description=Run task sync every 15 min  and on boot

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
