# .config/systemd/user/tasksync.service
[Unit]
Description=taskwarrior sync

[Service]
ExecStart=%h/.config/systemd/user/tasksync.sh
Type=oneshot

[Install]
WantedBy=default.target
