# .config/systemd/user/tasksync.service
[Unit]
Description=taskwarrior sync

[Service]
ExecStart=task sync
Type=oneshot

[Install]
WantedBy=default.target
