# .config/systemd/user/vdirsyncer.service
[Unit]
Description=vdirsyncer sync

[Service]
ExecStart=vdirsyncer sync
Type=oneshot

[Install]
WantedBy=default.target
