# .config/systemd/user/vdirsyncer.service
[Unit]
Description=vdirsyncer sync

[Service]
ExecStart=sh -c "vdirsyncer sync"
Type=oneshot

[Install]
WantedBy=default.target
