# .config/kanata/service/org.bferdinandy.kanata.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.bferdinandy.kanata</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/kanata</string>
        <string>-c</string>
        <string>/Users/bence.ferdinandy/.config/kanata/mac-entry.kbd</string>
        <string>-p</string>
        <string>7071</string>
    </array>

    <key>StandardOutPath</key>
    <string>/var/log/kanata.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/kanata.err</string>

    <!-- KeepAlive with PathState: launchd starts kanata only once the vhidd
         socket exists, and restarts it if it crashes (e.g. after sleep/wake).
         Also stops kanata if the socket disappears. -->
    <key>SessionCreate</key>
    <true/>

    <key>Nice</key>
    <integer>-10</integer>

    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/Library/Application Support/org.pqrs/tmp/rootonly/vhidd_server</key>
            <true/>
        </dict>
    </dict>
</dict>
</plist>
